# systemd/dictation.service
[Unit]
Description=dicti local dictation daemon (socket listener, recorder, insert)
After=graphical-session.target whisper-server.service
Wants=whisper-server.service

[Service]
Type=simple
# @SRC@ is substituted with the repo's src/ dir by install/04-install-daemon.sh
Environment=PYTHONPATH=@SRC@
ExecStart=/usr/bin/python3 -m dicti.daemon
Restart=on-failure
RestartSec=2
# Environment for ydotool, xclip, notify-send (set by graphical-session)
PassEnvironment=DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS YDOTOOL_SOCKET

[Install]
WantedBy=graphical-session.target
