# systemd/whisper-server.service
[Unit]
Description=whisper.cpp HTTP inference server (Vulkan, model hot in RAM)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Vulkan ACL on /dev/dri/renderD128 is applied by systemd-logind asynchronously
# with graphical-session.target. Probe until vulkaninfo sees the GPU (max 30s)
# so we never silently fall back to CPU.
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do vulkaninfo --summary 2>/dev/null | grep -q PHYSICAL_DEVICE_TYPE_INTEGRATED_GPU && exit 0; sleep 1; done; exit 1'
ExecStart=%h/opt/whisper.cpp/build/bin/whisper-server \
  --model %h/opt/whisper.cpp/models/ggml-medium-q5_0.bin \
  --host 127.0.0.1 \
  --port 8178 \
  --threads 4 \
  --suppress-nst \
  --inference-path /inference \
  --vad \
  --vad-model %h/opt/whisper.cpp/models/ggml-silero-v5.1.2.bin \
  --vad-threshold 0.25 \
  --vad-min-speech-duration-ms 0 \
  --vad-min-silence-duration-ms 100 \
  --vad-speech-pad-ms 500
# VAD drops silence so the streaming daemon never transcribes (and hallucinates on) it.
# threshold 0.25 + 500ms speech padding are deliberately generous: a tighter VAD trims quiet
# speech ONSETS, eating the first word after a pause on a low-output mic. See docs/v0.3-streaming.md.
Restart=on-failure
RestartSec=2
StartLimitBurst=5
StartLimitIntervalSec=60
# To switch to CPU fallback if Vulkan breaks: change build/ -> build-cpu/ above
# AND comment out the ExecStartPre line.

[Install]
WantedBy=graphical-session.target
